# usr/lib/udev/rules.d/96-adapter-notify.rules
# handle change of power source ac/bat, ignore input device batteries
ACTION=="change", SUBSYSTEM=="power_supply", KERNEL!="hidpp_battery*", RUN+="/usr/local/bin/sigblock -b bat"
